# data/openstargazer.service
[Unit]
Description=OpenStarGazer – Tobii Eye Tracker 5 Daemon
Documentation=https://github.com/1psconstructor/openstargazer
After=graphical-session.target

[Service]
Type=simple
# osg-setup install-service should set the correct path for ExecStart.
# The Environment line ensures ~/.local/bin is on PATH for pip-installed scripts.
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/bin/sh -c 'PATH="$HOME/.local/bin:$PATH" exec osg-daemon'
Restart=on-failure
RestartSec=3

[Install]
WantedBy=graphical-session.target
